# config/supervisor.conf
; Notes:
; priority=1 --> Lower priorities indicate programs that start first and shut down last
; killasgroup=true --> send kill signal to child processes too

[program:frappe-web]
command=/home/basawaraj/frappe-bench/env/bin/gunicorn -b 127.0.0.1:8000 -w 2 -t 120 frappe.app:application
priority=4
autostart=true
autorestart=true
stdout_logfile=/home/basawaraj/frappe-bench/logs/web.log
stderr_logfile=/home/basawaraj/frappe-bench/logs/web.error.log
user=basawaraj
directory=/home/basawaraj/frappe-bench/sites

[program:frappe-worker]
command=/home/basawaraj/frappe-bench/env/bin/python -m frappe.celery_app worker -n jobs@%%h --soft-time-limit 360 --time-limit 390 --loglevel INFO -Ofair
priority=4
autostart=true
autorestart=true
stdout_logfile=/home/basawaraj/frappe-bench/logs/worker.log
stderr_logfile=/home/basawaraj/frappe-bench/logs/worker.error.log
user=basawaraj
stopwaitsecs=400
directory=/home/basawaraj/frappe-bench/sites
killasgroup=true

[program:frappe-longjob-worker]
command=/home/basawaraj/frappe-bench/env/bin/python -m frappe.celery_app worker -n longjobs@%%h --soft-time-limit 1500 --time-limit 1530 --loglevel INFO
priority=2
autostart=true
autorestart=true
stdout_logfile=/home/basawaraj/frappe-bench/logs/worker.log
stderr_logfile=/home/basawaraj/frappe-bench/logs/worker.error.log
user=basawaraj
stopwaitsecs=1540
directory=/home/basawaraj/frappe-bench/sites
killasgroup=true

[program:frappe-async-worker]
command=/home/basawaraj/frappe-bench/env/bin/python -m frappe.celery_app worker -n async@%%h --soft-time-limit 1500 --time-limit 1530 --loglevel INFO
priority=2
autostart=true
autorestart=true
stdout_logfile=/home/basawaraj/frappe-bench/logs/worker.log
stderr_logfile=/home/basawaraj/frappe-bench/logs/worker.error.log
user=basawaraj
stopwaitsecs=1540
directory=/home/basawaraj/frappe-bench/sites
killasgroup=true

[program:frappe-workerbeat]
command=/home/basawaraj/frappe-bench/env/bin/python -m frappe.celery_app beat -s beat.schedule
priority=3
autostart=true
autorestart=true
stdout_logfile=/home/basawaraj/frappe-bench/logs/workerbeat.log
stderr_logfile=/home/basawaraj/frappe-bench/logs/workerbeat.error.log
user=basawaraj
directory=/home/basawaraj/frappe-bench/sites


[program:redis-cache]
command=/usr/bin/redis-server /home/basawaraj/frappe-bench/config/redis_cache.conf
priority=1
autostart=true
autorestart=true
stdout_logfile=/home/basawaraj/frappe-bench/logs/redis-cache.log
stderr_logfile=/home/basawaraj/frappe-bench/logs/redis-cache.error.log
user=basawaraj
directory=/home/basawaraj/frappe-bench/sites

[program:redis-async-broker]
command=/usr/bin/redis-server /home/basawaraj/frappe-bench/config/redis_async_broker.conf
priority=1
autostart=true
autorestart=true
stdout_logfile=/home/basawaraj/frappe-bench/logs/redis-async-broker.log
stderr_logfile=/home/basawaraj/frappe-bench/logs/redis-async-broker.error.log
user=basawaraj
directory=/home/basawaraj/frappe-bench/sites

[program:node-socketio]
command=/usr/bin/nodejs /home/basawaraj/frappe-bench/apps/frappe/socketio.js
priority=4
autostart=true
autorestart=true
stdout_logfile=/home/basawaraj/frappe-bench/logs/node-socketio.log
stderr_logfile=/home/basawaraj/frappe-bench/logs/node-socketio.error.log
user=basawaraj
directory=/home/basawaraj/frappe-bench/sites


[group:frappe]
programs=frappe-web,frappe-worker,frappe-workerbeat